# op_metrics/rootfs/etc/services.d/op_metrics/finish
#!/usr/bin/with-contenv bashio
bashio::log.info "op-metrics exited with code ${1}"
